# configs/systemd/pdns.service.d/dnslb.conf
[Service]
Restart=always

[Unit]
Requisite=dnslb.service
